# kautoswitch/packaging/systemd/kautoswitch.service
[Unit]
Description=KAutoSwitch — keyboard layout auto-corrector daemon
Documentation=https://github.com/kautoswitch/kautoswitch
After=graphical-session.target
PartOf=graphical-session.target

[Service]
Type=simple
ExecStart=/usr/bin/kautoswitch-daemon
Restart=on-failure
RestartSec=1s
# Avoid restart storms
StartLimitIntervalSec=30
StartLimitBurst=5
# Environment: DISPLAY is inherited from user session on KDE
# If not, uncomment: Environment=DISPLAY=:0

[Install]
WantedBy=graphical-session.target
